# deploy/systemd/ssnp-probe.service
[Unit]
Description=SSNP Probe Worker
After=network-online.target ssnp-portal.service
Wants=network-online.target

[Service]
Type=simple
User=ssnp-probe
Group=ssnp-probe
WorkingDirectory=/opt/ssnp-probe
ExecStart=/opt/ssnp-probe/bin/probe-worker --config /etc/ssnp-probe/config.yaml run
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
